# .config/systemd/user/update-enabled-units.timer
[Unit]
Description=Run update-enabled-units daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
